# templates/etc/profile.d/influxdb.sh.j2
export PATH=$PATH:{{influxdb_install_parent_dir}}/influxdb/usr/bin